# update-deployment.timer
[Unit]
Description=Timer para revision de cambios

[Timer]
OnCalendar=*-*-* *:00:00
Unit=update-deployment.service

[Install]
WantedBy=timers.target
