# update-deployment.timer
[Unit]
Description=Timer para revisión de cambios

[Timer]
OnUnitActiveSec=30
Unit=update-deployment.service

[Install]
WantedBy=timers.target
